# machines/chef/cookbooks/coursys_deployed/files/default/rc.local
#!/bin/sh -e

su coursys -c supervisord
